# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for pmpleka.cloudapp.net"

start on net-device-up
stop on shutdown

respawn

setuid elspeth
chdir /home/elspeth/sites/pmpleka.cloudapp.net/source

exec ../virtualenv/bin/gunicorn \
    --bind unix:/tmp/pmpleka.cloudapp.net.socket \
    superlists.wsgi:application
